# solocloud.service
[Unit]
Description=SoloCloud Personal Cloud Storage
After=network.target

[Service]
Type=notify
User=SoloCloud
Group=SoloCloud
WorkingDirectory=/opt/SoloCloud
Environment=FLASK_ENV=production
Environment=FLASK_APP=app.py
ExecStart=/opt/SoloCloud/venv/bin/gunicorn --config gunicorn.conf.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/SoloCloud /var/log/SoloCloud /var/run/SoloCloud
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
